# Unicorn_Ruby/supervisor_unicorn.conf
[program:unicorn]

command=bash /home/sergiy/start_unicorn.sh

autostart=true
autorestart=true

stderr_logfile=/var/log/unicorn.err.log
stdout_logfile=/var/log/unicorn.out.log

user=root
